# hdSleep.timer
[Unit]
Description=Run hdSleep each 2min and 5min after boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
