# hdSleep.service

[Unit]
Description=Run hdSleep each 2min and 5min after boot

[Service]
Type=oneshot
ExecStart=/bin/sh /home/patrick/Documentos/hdSleep.sh
